# packages/discord-voicebot/discord_voicebot-0.1.0.tar.gz/discord_voicebot-0.1.0/data/voicebot@.service
[Unit]
Description=VoiceBot instance (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/voicebot --token=%i
# fallback if ExecStart lacks --token
EnvironmentFile=/etc/voicebot/%i.env
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
